# system_files/shared/usr/lib/systemd/user/xdg-desktop-portal.service.d/30-after-keyring.conf
[Unit]
# Ensure gnome-keyring-daemon is running before the portal starts,
# so the Secret portal backend (org.freedesktop.portal.Secret) is
# available when xdg-desktop-portal probes for backends.
# Without this, the portal may start before the keyring and cache
# stale backend discovery, causing apps like Pika Backup to fail
# with "background process inactive" or keyring access errors.
After=gnome-keyring-daemon.service
Wants=gnome-keyring-daemon.service
