# system_files/kinoite/usr/etc/profile.d/ksshaskpass.sh
#!/bin/sh

SSH_ASKPASS=/usr/bin/ksshaskpass
export SSH_ASKPASS
