# usr/etc/profile.d/brew.sh
[[ $(command -v brew >/dev/null 2>&1) && $- == *i* ]] && eval "$(/home/linuxbrew/.linuxbrew/bin/brew shellenv)"