# system_files/shared/usr/lib/systemd/user/bluefin-dynamic-wallpaper.timer
[Unit]
Description=Run Bluefin Dynamic Wallpaper daily
Wants=graphical-session.target
After=graphical-session.target

[Timer]
OnBootSec=1min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=graphical-session.target
