# usr/lib/systemd/user/bluefin-cli.target
[Unit]
Description=Bluefin-CLI target for bluefin-cli quadlet

[Install]
WantedBy=default.target