# usr/lib/systemd/user/ublue-user-setup.service
[Unit]
Description=Configure system for current user

[Service]
Type=simple
ExecStart=/usr/bin/ublue-user-setup

[Install]
WantedBy=default.target
