# etc/systemd/system/nix-daemon.service.d/override.conf
[Service]
Environment="NIX_SSL_CERT_FILE=/etc/ssl/certs/ca-bundle.crt"