# system_files/shared/etc/profile.d/91-bluefin-aliases.sh
# Bluefin convenience aliases
# ramalama shorthand - https://github.com/containers/ramalama
if command -v ramalama &>/dev/null; then
    alias rl='ramalama'
fi
