# system_files/shared/etc/profile.d/askpass.sh
#!/usr/bin/bash

SUDO_ASKPASS='/usr/libexec/openssh/gnome-ssh-askpass'
export SUDO_ASKPASS
