# system_files/dx/usr/lib/systemd/system/bluefin-dx-brew.service
[Unit]
Description=Run bluefin-dx-brew once on user login

[Service]
Type=oneshot
ExecStart=/usr/libexec/bluefin-dx-brew

[Install]
WantedBy=default.target
